# docker/docker.mongodb.service
[Unit]
Description=MongoDB Docker Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker exec %n stop
ExecStartPre=-/usr/bin/docker rm %n
ExecStartPre=/usr/bin/docker pull mongo:4.4.8
ExecStart=/usr/bin/docker run --rm --name %n \
    -v /opt/mongodb:/data/db \
    -p 27017:27017 \
    mongo:4.4.8

[Install]
WantedBy=default.target
