# pets/rox_stuff/documents_backup/documents-backup.timer
[Unit]
Description=Backup documents daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
